# extra/spec/foreman-proxy-abrt-send.cron
# Send collected ABRT reports once every 10 minutes
*/30 * * * * foreman-proxy /usr/share/foreman-proxy/bin/smart-proxy-abrt-send
